# ansible/roles/monitoring_server/templates/monitoring.target
# {{ ansible_managed }}

[Unit]
Description=Monitoring functionality
After=network.target

[Install]
WantedBy=default.target
